# contrib/supervisor/pybossa-analyst.conf.tmpl
[program:pybossa-analyst]
command=/home/pybossa-analyst/env/bin/uwsgi /home/pybossa-analyst/pybossa-analyst.ini
directory=/home/pybossa-analyst
autostart=true
autorestart=true
user=pybossa