# contrib/supervisor/rq-worker.conf.tmpl
[program:rq-worker]
command=/home/pybossa-analyst/env/bin/python app_context_rqworker.py pybossa-analyst
directory=/home/pybossa-analyst
autostart=true
autorestart=true
user=pybossa
priority=997